# launchd/net.consolaktif.discord.launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.consolaktif.discord.launcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>
timeout 30 bash -c 'until lsof -i :${CD_PROXY_PORT:-8080} >/dev/null 2>&1; do sleep 0.2; done'
# Discord’u başlat
exec /Applications/Discord.app/Contents/MacOS/Discord --proxy-server=http://127.0.0.1:${CD_PROXY_PORT:-8080} &
# 5 sn sonra proxy bağlantısını kontrol et, yoksa bildirim at
(sleep 5 && netstat -an | grep -q ":${CD_PROXY_PORT:-8080}.*ESTABLISHED" || osascript -e 'display notification "Discord proxy’ye bağlanamadı" with title "SplitWire"') &
        </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>HardResourceLimits</key>
    <dict>
        <key>CPU</key><integer>60</integer>
        <key>Memory</key><integer>268435456</integer>
    </dict>
</dict>
</plist>